# docker/assets/etc/service/nginx/run
#!/bin/bash
exec 2>&1
source /usr/local/etc/envvars
exec /usr/sbin/nginx
