# docker/assets/etc/service/php-fpm/run
#!/bin/bash
exec 2>&1
source /usr/local/etc/envvars
php-fpm-env > /usr/local/etc/php-fpm-env.conf
exec /usr/local/sbin/php-fpm --nodaemonize
